# RemoteCoOp/service/macos/com.macforce-now.remote-coop.panel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.macforce-now.remote-coop.panel</string>
  <key>ProgramArguments</key>
  <array>
    <string>__NODE__</string>
    <string>RemoteCoOp/panel/control-panel.mjs</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__</string>
  <key>UserName</key>
  <string>__SERVICE_USER__</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>MACFORCE_NOW_REMOTE_COOP_PANEL_BIND_HOST</key>
    <string>0.0.0.0</string>
    <key>MACFORCE_NOW_REMOTE_COOP_PANEL_PORT</key>
    <string>__PANEL_PORT__</string>
    <key>MACFORCE_NOW_REMOTE_COOP_PUBLIC_HOST</key>
    <string>198.12.95.48</string>
    <key>MACFORCE_NOW_REMOTE_COOP_PORT</key>
    <string>__BROKER_PORT__</string>
    <key>MACFORCE_NOW_REMOTE_COOP_PORT_ALTERNATES</key>
    <string>__BROKER_ALTERNATES__</string>
    <key>MACFORCE_NOW_REMOTE_COOP_TURN_PORT</key>
    <string>__TURN_PORT__</string>
    <key>MACFORCE_NOW_REMOTE_COOP_TURN_TLS_PORT</key>
    <string>32443</string>
    <key>MACFORCE_NOW_REMOTE_COOP_TURN_MIN_PORT</key>
    <string>__TURN_MIN_PORT__</string>
    <key>MACFORCE_NOW_REMOTE_COOP_TURN_MAX_PORT</key>
    <string>__TURN_MAX_PORT__</string>
    <key>MACFORCE_NOW_REMOTE_COOP_AUTOSTART</key>
    <string>1</string>
    <key>MACFORCE_NOW_REMOTE_COOP_PANEL_UPDATE_AUTOMATIC</key>
    <string>0</string>
    <key>MACFORCE_NOW_REMOTE_COOP_TURN_SHARED_SECRET</key>
    <string>__TURN_SHARED_SECRET__</string>
    <key>MACFORCE_NOW_REMOTE_COOP_INVITE_SECRET</key>
    <string>__INVITE_SECRET__</string>
  </dict>
  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/var/log/macforce-now-remote-coop-panel.log</string>
  <key>StandardErrorPath</key>
  <string>/var/log/macforce-now-remote-coop-panel.err.log</string>
</dict>
</plist>
